# lte/setup-lte.service
[Unit] 
Description=Set up LTE modem in the dumbest possible way 

[Service]
ExecStart=python /rno-g/bin/setup-lte.py 
